# deploy/ubuntu/real-estate-watcher.service
[Unit]
Description=real-estate-watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=realestate
Group=realestate
WorkingDirectory=/opt/real-estate-watcher
EnvironmentFile=/etc/real-estate-watcher/real-estate-watcher.env
ExecStart=/usr/bin/env node /opt/real-estate-watcher/dist/src/index.js run-service
Restart=on-failure
RestartSec=5
StateDirectory=real-estate-watcher
CacheDirectory=real-estate-watcher
LogsDirectory=real-estate-watcher
TimeoutStopSec=30
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=real-estate-watcher

[Install]
WantedBy=multi-user.target
